# systemd/agent-alert-monitor-ingest.service
[Unit]
Description=application alert Telegram intake
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
UMask=0077
WorkingDirectory=__ALERT_MONITOR_ROOT_SYSTEMD__
EnvironmentFile=-__ALERT_MONITOR_ENV_SYSTEMD__
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__ALERT_MONITOR_BIN_SYSTEMD__ --config __ALERT_MONITOR_CONFIG_SYSTEMD__ listen
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
