# systemd/agent-alert-monitor-watchdog.service
[Unit]
Description=application alert watchdog no-silence check

[Service]
Type=oneshot
UMask=0077
WorkingDirectory=__ALERT_MONITOR_ROOT_SYSTEMD__
EnvironmentFile=-__ALERT_MONITOR_ENV_SYSTEMD__
ExecStart=__ALERT_MONITOR_BIN_SYSTEMD__ --config __ALERT_MONITOR_CONFIG_SYSTEMD__ watchdog-due --send-telegram
