# systemd/agent-alert-monitor-health.service
[Unit]
Description=application alert SQS health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
UMask=0077
WorkingDirectory=__ALERT_MONITOR_ROOT_SYSTEMD__
EnvironmentFile=-__ALERT_MONITOR_ENV_SYSTEMD__
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=ALERT_MONITOR_SQS_SOURCE=sample-api-prod-alerts
ExecStart=__ALERT_MONITOR_BIN_SYSTEMD__ --config __ALERT_MONITOR_CONFIG_SYSTEMD__ health --source ${ALERT_MONITOR_SQS_SOURCE} --json
